# docker/config/supervisord.conf
# ABOUTME: Supervisord config — manages code-server and LiteLLM proxy processes.
# ABOUTME: code-server runs via its upstream entrypoint; LiteLLM proxies OpenAI.

[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:code-server]
command=/usr/bin/entrypoint.sh --bind-addr 0.0.0.0:8080 .
autostart=true
autorestart=true
priority=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# Cisco TLS inspection proxy: uncomment the line below if MCP servers fail with
# TLS errors. Node.js does not use the system cert store, so the CA injection in
# the Containerfile alone is not sufficient for Node processes.
# environment=NODE_TLS_REJECT_UNAUTHORIZED="0"

[program:litellm]
command=litellm --config /root/config/litellm-config.yaml --host 127.0.0.1 --port 4000
autostart=true
autorestart=true
priority=2
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
